# assets/templates/browserbase/cua-local/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:cua-server]
command=/app/cua-local-server/cua-local-server-linux-x64
directory=/app/cua-local-server
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/cua-server.err.log
stdout_logfile=/var/log/supervisor/cua-server.out.log
environment=CUA_SERVER_PORT="3001",CUA_SERVER_HOST="0.0.0.0"

[program:nextjs-app]
command=npm run start
directory=/app/nextjs-app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nextjs-app.err.log
stdout_logfile=/var/log/supervisor/nextjs-app.out.log
environment=PORT="3000"
